# agent/packaging/ubuntu/kontena-weave/etc/init/kontena-weave-helper.conf
description "Kontena weave helper"

start on started docker
stop on stopping docker

respawn

script
  /usr/local/bin/weave-helper
end script
