# agent/packaging/ubuntu_xenial/kontena-agent/etc/systemd/system/docker.service.d/50-kontena.conf
[Service]
Environment='DOCKER_OPTS=--insecure-registry="10.81.0.0/16"'
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// ${DOCKER_OPTS}
