# agent/packaging/ubuntu_xenial/kontena-agent/lib/systemd/system/docker.service.d/50-kontena.conf
[Service]
Environment='DOCKER_OPTS=--insecure-registry="10.81.0.0/19"'
