# agent/packaging/ubuntu_xenial/kontena-agent/DEBIAN/prerm
#!/bin/sh

systemctl stop kontena-agent || true

systemctl disable kontena-agent || true
