# agent/packaging/ubuntu_xenial/kontena-agent/DEBIAN/prerm
#!/bin/sh

/opt/kontena/kontena-resolvconf.sh stop

systemctl stop kontena-agent || true

systemctl disable kontena-agent || true
