# server/packaging/ubuntu_xenial/kontena-server/lib/systemd/system/kontena-server.service.d/kontena-version.conf
[Service]
Environment=KONTENA_VERSION={{VERSION}}
